# templates/docker.socket.j2
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -H tcp://0.0.0.0:{{docker_remote_port}} -H unix://var/run/docker.sock
